# deploy/calendar-backend.service
[Unit]
Description=Calendar Kiosk Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/calendar-kiosk/backend
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/opt/calendar-kiosk/backend/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
